# lib/44-midipipes-bt.rules
ACTION=="add|remove", SUBSYSTEM=="bluetooth", RUN+="/usr/bin/python {MIDIPI_USER}/bin/midi.py"
