# lib/midipipes-web.service
[Unit]
Description=MIDI Pipes Web

[Service]
Type=simple
ExecStart="{MIDIPI_USER}/bin/web.py"

[Install]
WantedBy=multi-user.target
